# systemd/llm-ocr-middleware.service
[Unit]
Description=LLM Gateway — OCR middleware (:8900)
After=network.target llm-litellm.service

[Service]
Type=simple
User=derek
Group=derek
WorkingDirectory=/home/derek/Projects/llm-gateway
EnvironmentFile=-/home/derek/Projects/llm-gateway/gateway.env
ExecStart=/home/derek/Projects/llm-gateway/run_ocr_middleware.sh
Restart=on-failure
RestartSec=10
TimeoutStartSec=60
TimeoutStopSec=15
StandardOutput=append:/home/derek/Projects/llm-gateway/logs/ocr_middleware.log
StandardError=append:/home/derek/Projects/llm-gateway/logs/ocr_middleware.log

[Install]
WantedBy=multi-user.target
